# ansible/roles/common_cli/files/systemd/wpa_supplicant@.service.d/override.conf
# Managed by Ansible - do not edit manually
# Make wireless device optional to avoid blocking boot
[Unit]
# Remove the strict Requires dependency on the device
# The service will still wait for the device due to After=, but won't fail/block if it's slow
Requires=
After=sys-subsystem-net-devices-%i.device
