# ansible/roles/common_cli/files/udev/80-wireless-wpa.rules
# Managed by Ansible - do not edit manually
# Automatically enable wpa_supplicant for wireless interfaces
ACTION=="add", SUBSYSTEM=="net", KERNEL=="wl*", RUN+="systemctl --no-block start wpa_supplicant@%k.service"
