# ansible/roles/common_cli/files/systemd/systemd-networkd-wait-online.service.d/timeout.conf
# Managed by Ansible - do not edit manually
# Reduce timeout so boot doesn't hang waiting for slow interfaces
[Service]
TimeoutStartSec=30
