# ansible/roles/common_cli/files/systemd/systemd-networkd-wait-online.service.d/timeout.conf
# Managed by Ansible - do not edit manually
# Note: This drop-in is not currently deployed by Ansible tasks. The playbook
# masks the systemd-networkd-wait-online service entirely instead of applying
# a timeout. In that case, this timeout has no effect and this file is kept
# for optional or future use.
#
# Reduce timeout so boot doesn't hang waiting for slow interfaces
[Service]
TimeoutStartSec=30
