# ansible/roles/clawdbot/templates/clawdbot.service.j2
[Unit]
Description=Clawdbot Gateway
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User={{ clawdbot_user }}
Group={{ clawdbot_user }}
WorkingDirectory={{ clawdbot_workspace_path }}

# fnm environment setup
Environment="PATH=/home/{{ clawdbot_user }}/.local/share/fnm/aliases/default/bin:/usr/local/bin:/usr/bin:/bin"
Environment="FNM_DIR=/home/{{ clawdbot_user }}/.local/share/fnm"

# 1Password service account token loaded from root-only file (mode 0600)
EnvironmentFile=-/etc/clawdbot/op.env

# Gateway command
ExecStart=/home/{{ clawdbot_user }}/.local/share/fnm/aliases/default/bin/clawdbot gateway --bind {{ clawdbot_bind }} {{ clawdbot_gateway_extra_args }}

# Restart policy
Restart=always
RestartSec=10

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths={{ clawdbot_workspace_path }} {{ clawdbot_config_dir }}

[Install]
WantedBy=multi-user.target
